# systemd/node_service.service
[Unit]
Description=Crypto API Node.js Gateway Service
After=network.target python_service.service
Requires=python_service.service

[Service]
User=root
WorkingDirectory=/root/crypto-api
EnvironmentFile=/root/crypto-api/.env
ExecStart=/usr/bin/npx tsx server/index.ts
Restart=always
RestartSec=5
TimeoutStartSec=180
StandardOutput=journal
StandardError=journal
Environment="NODE_ENV=production"
Environment="NODE_OPTIONS=--max-old-space-size=2048"

[Install]
WantedBy=multi-user.target
