# systemd/python_service.service
[Unit]
Description=Crypto API Python Service (CoinGlass FastAPI)
Documentation=https://github.com/rcz87/crypto-api
After=network-online.target postgresql.service redis.service
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=3

[Service]
Type=notify
User=root
Group=root
WorkingDirectory=/root/crypto-api/coinglass-system

# Environment Configuration
EnvironmentFile=/root/crypto-api/.env
Environment="PYTHONUNBUFFERED=1"
Environment="PYTHONPATH=/root/crypto-api/coinglass-system"
Environment="LOG_LEVEL=INFO"

# Service Execution
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/python3 -m uvicorn app.main:app \
    --host 0.0.0.0 \
    --port 8000 \
    --workers 4 \
    --log-level info \
    --access-log \
    --use-colors

# Graceful Shutdown
ExecReload=/bin/kill -HUP $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30
RestartSec=5

# Auto-restart Configuration
Restart=always
RestartSec=10

# Security Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/root/crypto-api/coinglass-system /var/log/crypto-api
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource Limits
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=python_service

# Watchdog (auto-restart if service hangs)
WatchdogSec=60

[Install]
WantedBy=multi-user.target
